# crontab
*/5 * * * * /data/script.sh > /dev/null 2>&1
